# support/build/rpmbuild/BUILD/flower-x.x/.defaults/flower
#!/bin/bash
#
# flower
#
# chkconfig: 2345 69 16
# description:  Start up the Flower engine.

# Source function library.
. /etc/init.d/functions


RETVAL=$?
RUN_AS_USER=flower
ExecStart=PUT_START_STRING_HERE

case "$1" in
 start)
          echo $"Starting Flower"
          /bin/su - $RUN_AS_USER $ExecStart
        ;;
 stop)
          echo $"Stopping Flower"
          /bin/su - $RUN_AS_USER killall flower
        ;;
 *)
        echo $"Usage: $0 {start|stop}"
        exit 1
        ;;
esac

exit $RETVAL
